# gunicorn/gunicorn.socket
[Unit]
Description=gunicorn socket
[Socket]
ListenStream=/home/ubuntu/project-pipeline-aws/aws_cicd.sock
[Install]
WantedBy=sockets.target
